# sys/countitperiodic.timer
[Unit]
Description=countit Periodic program

[Timer]
OnBootSec=5min
OnUnitActiveSec=10min

[Install]
WantedBy = multi-user.target